# osx/init/Library/LaunchDaemons/com.discoveryd.fix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.fix.local</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/sbin/discoveryutil</string>
            <string>activedirectory</string>
            <string>yes</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>ServiceDescription</key>
        <string>startup service for domain.local discoveryd fix</string>
    </dict>
</plist>